# supervisor.conf
[program:workbench]
command=make run
directory=/home/vagrant
user=vagrant
autostart=true
